# crontab
0 */1 * * * /usr/sbin/ntpdate -u 172.17.4.44 > /dev/null
0 */1 * * * /sbin/hwclock -w
